# infra/systemd/library-frontend.service
[Unit]
Description=Library Platform Next.js Frontend
After=network.target

[Service]
Type=simple
User=fg
Group=fg
WorkingDirectory=/home/fg/library-platform/frontend
EnvironmentFile=/home/fg/library-platform/frontend/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
